# apps/core/email_admin/setup/stats/postfix_stats.cron
01 * * * * root /usr/bin/php -q /usr/local/elastix/postfix_stats.php
