# apps/core/security/setup/etc/cron.d/elastix-portknock.cron
0 * * * * root /usr/bin/elastix-portknock-cleanup